# prod_deploy/supervisord.conf
; Wosobo multi-process supervisor
; Shared secrets/URLs come from container env (docker-compose / .env).

[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:phone-server]
directory=/opt/wosobo/phone-server
command=node src/index.js
autostart=true
autorestart=true
startsecs=2
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="%(ENV_PHONE_SERVER_PORT)s",WS_PORT="%(ENV_PHONE_SERVER_WS_PORT)s",MONGODB_URI="%(ENV_MONGODB_URI)s",CORS_ORIGIN="%(ENV_CORS_ORIGIN)s",JANUS_WS_URL="%(ENV_JANUS_WS_URL)s",INTERNAL_TOKEN="%(ENV_INTERNAL_TOKEN)s",ABSENT_ANNOUNCE_FILE="%(ENV_ABSENT_ANNOUNCE_FILE)s",ABSENT_ANNOUNCE_MAX_SEC="%(ENV_ABSENT_ANNOUNCE_MAX_SEC)s",CALL_CDR_TTL_SEC="%(ENV_CALL_CDR_TTL_SEC)s"

[program:manage-api]
directory=/opt/wosobo/manage-api
command=node src/index.js
autostart=true
autorestart=true
startsecs=2
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="%(ENV_MANAGE_API_PORT)s",MONGODB_URI="%(ENV_MONGODB_URI)s",CORS_ORIGIN="%(ENV_CORS_ORIGIN)s",MANAGE_API_TOKEN="%(ENV_MANAGE_API_TOKEN)s",OPENAPI_SERVER_URL="%(ENV_OPENAPI_SERVER_URL)s",PHONE_SERVER_URL="%(ENV_PHONE_SERVER_URL)s",INTERNAL_TOKEN="%(ENV_INTERNAL_TOKEN)s"

[program:monitor]
directory=/opt/wosobo/monitor
command=node src/index.js
autostart=true
autorestart=true
startsecs=1
stopwaitsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="%(ENV_MONITOR_PORT)s",JANUS_ADMIN_URL="%(ENV_JANUS_ADMIN_URL)s",JANUS_ADMIN_SECRET="%(ENV_JANUS_ADMIN_SECRET)s"

[program:softphone-demo]
directory=/opt/wosobo/softphone-demo
command=node src/index.js
autostart=true
autorestart=true
startsecs=1
stopwaitsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="%(ENV_SOFTPHONE_DEMO_PORT)s",MANAGE_API_URL="%(ENV_MANAGE_API_URL)s",MANAGE_API_TOKEN="%(ENV_MANAGE_API_TOKEN)s",CORS_ORIGIN="%(ENV_CORS_ORIGIN)s"

[program:static-server]
directory=/opt/wosobo/static-server
command=node src/index.js
autostart=true
autorestart=true
startsecs=1
stopwaitsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PORT="%(ENV_STATIC_PORT)s",STATIC_ROOT="%(ENV_STATIC_ROOT)s"
